# backend/scripts/systemd/autophon-worker.service
[Unit]
Description=AutoPhon Alignment Worker
Documentation=file:///path/to/your/app/backend/scripts/README.md
After=network.target postgresql.service
Wants=network.target
Requires=postgresql.service

[Service]
Type=simple
User=nordalign
Group=nordalign
WorkingDirectory=/home/nordalign/new-autophon/backend
ExecStart=/usr/bin/python3 /home/nordalign/new-autophon/backend/scripts/alignment_worker.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

# Environment
Environment=PYTHONPATH=/home/nordalign/new-autophon/backend
EnvironmentFile=/home/nordalign/new-autophon/backend/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/nordalign/new-autophon/uploads /home/nordalign/new-autophon/logs /tmp
PrivateTmp=true

# Resource limits
MemoryLimit=2G
TasksMax=50
LimitNOFILE=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autophon-worker

[Install]
WantedBy=multi-user.target