# backend/scripts/systemd/autophon-upload-worker.service
[Unit]
Description=AutoPhon Real-time Upload Worker
Documentation=file:///home/nordalign/new-autophon/backend/scripts/README.md
After=network.target postgresql.service
Wants=network.target
Requires=postgresql.service

[Service]
Type=simple
User=nordalign
Group=nordalign
WorkingDirectory=/home/nordalign/new-autophon/backend
ExecStart=/usr/bin/python3 /home/nordalign/new-autophon/backend/scripts/upload_worker.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=2
KillMode=mixed
TimeoutStopSec=15

# Environment
Environment=PYTHONPATH=/home/nordalign/new-autophon/backend
EnvironmentFile=/home/nordalign/new-autophon/backend/.env

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/nordalign/new-autophon/uploads /home/nordalign/new-autophon/logs /tmp
PrivateTmp=true

# Resource limits (lighter than alignment worker)
MemoryLimit=1G
TasksMax=25
LimitNOFILE=2048

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autophon-upload-worker

# Real-time processing optimization
Nice=-5
IOSchedulingClass=1
IOSchedulingPriority=4

[Install]
WantedBy=multi-user.target